# com.nuthouse.stagepad-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nuthouse.stagepad-bridge</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>-m</string>
        <string>bridge.main</string>
    </array>

    <key>WorkingDirectory</key>
    <string>BRIDGE_DIR_PLACEHOLDER</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>BRIDGE_DIR_PLACEHOLDER</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>LOG_DIR_PLACEHOLDER/bridge.log</string>
    <key>StandardErrorPath</key>
    <string>LOG_DIR_PLACEHOLDER/bridge.log</string>

    <key>ThrottleInterval</key>
    <integer>3</integer>
</dict>
</plist>
