# com.nuthouse.stagepad-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nuthouse.stagepad-bridge</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>SCRIPTS_DIR_PLACEHOLDER/start-stagepad-bridge.sh</string>
    </array>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>LOG_DIR_PLACEHOLDER/bridge.log</string>
    <key>StandardErrorPath</key>
    <string>LOG_DIR_PLACEHOLDER/bridge.log</string>

    <key>ThrottleInterval</key>
    <integer>3</integer>
</dict>
</plist>
